# crontab/notify.sh
#!/bin/sh
/www/server/php/71/bin/php /home/wwwroot/erc20/public/index.php Shell_Orders/notifyMerch
